# Deploy/supervisord/run.d/nuxt_web.ini
[program:nuxt_web]
command=npm run start
directory=/home/app/Web
autostart=true
autorestart=unexpected
exitcodes=0
priority = 1